# nuxt-invoices.service
[Unit]
Description=Nuxt SSR (invoices)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/invoices-test/invoices.com
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=3000
Environment=NUXT_API_BASE=http://127.0.0.1/api
Environment=NUXT_PUBLIC_API_BASE=/api
ExecStart=/usr/local/bin/node /var/www/invoices-test/invoices.com/.output/server/index.mjs
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
